# coreos/roles/rwky_net/files/clean-docker.service
[Unit]
Description=cleans docker images

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "docker images | grep 'none' | xargs docker rmi" 

