# coreos/roles/rwky_net/files/clean-docker.timer
[Unit]
Description=clean docker

[Timer]
OnBootSec=15m
OnUnitInactiveSec=15m
Persistent=true     
 
[Install]
WantedBy=timers.target
